# tools/systemd/openmmo-server.service
[Unit]
Description=OpenMMO game server
Documentation=https://github.com/Julian-adv/OpenMMO
After=network-online.target
Wants=network-online.target
# Don't give up if the server crash-loops during a bad deploy.
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
# The binary resolves data/, data/terrain/ and data/npc_token relative to CWD.
WorkingDirectory=/home/ubuntu/work/OnlineRPG
EnvironmentFile=-/etc/openmmo/server.env
ExecStart=/home/ubuntu/work/OnlineRPG/target/release/onlinerpg-server

Restart=always
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=openmmo

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
